# AdtClaw/adt_openclaw/deploy/systemd/adt-openclaw-provider.service
[Unit]
Description=OpenClaw Local Provider
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=openclaw
Group=openclaw
WorkingDirectory=/var/lib/openclaw
EnvironmentFile=-/etc/openclaw/openclaw.env
Environment=OPENCLAW_BIND_ADDR=127.0.0.1:11434
ExecStart=/usr/local/bin/adt_openclaw_provider_launcher
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
