# AdtClaw/adt_openclaw/deploy/systemd/adt-openclaw-healthcheck.service
[Unit]
Description=OpenClaw Provider Healthcheck
After=adt-openclaw-provider.service
Requires=adt-openclaw-provider.service

[Service]
Type=oneshot
EnvironmentFile=-/etc/openclaw/openclaw.env
ExecStart=/usr/local/bin/adt_openclaw_healthcheck wait
